# systemd/skeleton-fpm-command.service
[Unit]
Description=Skeleton FPM - Each 20 min

[Service]
Type=oneshot
ExecStart=/opt/skeleton/systemd/scripts/skeleton-fpm-command.sh
